# etc/crontab/a2b-lax.sh
#!/bin/bash
export XAUTHORITY=/home/thoth/.Xauthority
export DISPLAY=:0
kdialog --msgbox="Visit go/a2b-lax"
